# etc/bb.d/rc.local
#!/lib/bb/sh
# -*- mode: sh; -*-

case "$1" in
  start)
    echo "Hello from $0"
  ;;
  stop)
    echo "Goodbye from $0"
  ;;
esac

